# src/root/etc/inittab
::sysinit:/bin/getty 38400 tty1
::restart:/sbin/init
::ctrlaltdel:/sbin/reboot
::once:cat /etc/motd
::respawn:/bin/cttyhack /bin/sh
